# deploy/scanner.service
# systemd service unit for Repository Scanner (on-premise)
#
# Install:
#   sudo cp deploy/scanner.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now repo-scanner
#
# Logs:
#   journalctl -u repo-scanner -f
#   docker compose -f /opt/repo-scanner/docker-compose.yml logs -f

[Unit]
Description=Repository Scanner
Documentation=file:///opt/repo-scanner/tutorials/deploying-to-rocky-linux.md
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes

WorkingDirectory=/opt/repo-scanner

# Start the full stack (Nginx + app) via Compose
ExecStart=/usr/bin/docker compose up -d

# Give in-flight scans 30 s to finish before stopping
ExecStop=/usr/bin/docker compose stop --timeout 30

# If the service is restarted, bring the stack down cleanly first
ExecStopPost=/usr/bin/docker compose down

Restart=on-failure
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
